# supervisor.conf
[supervisord]
logfile=/tmp/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid

[program: vscode-server]
command = /usr/bin/vscode.sh
autostart = true
autorestart = true

